# v1_1_5/dockerfiles/build/zenbership/configs/supervisor/supervisord.conf
[inet_http_server]           ; inet (TCP) server disabled by default
port = 44100                 ; (ip_address:port specifier, *:port for all iface)
username = user              ; (default is no username (open server))
password = anynewpassword    ; (default is no password (open server))

[supervisorctl]
username = dummy
password = anynewpassword

[supervisord]
nodaemon = true

[program:php-fpm]
command = php-fpm -F
autostart = true
autorestart = true
;stdout_logfile=/dev/stdout
;stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/php.err.log
stderr_logfile_maxbytes=5MB
stdout_logfile=/var/log/supervisor/php.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=7
redirect_stderr=true

[program:cron]
autorestart = false
command = cron -f
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[include]
files = supervisor/conf.d/*.conf